# infra/launchd/com.sprout.cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sprout.cloudflared</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/cloudflared</string>
    <string>--config</string>
    <string>/Users/waylantrieu/repos/sprout/infra/cloudflared/config.yml</string>
    <string>tunnel</string>
    <string>run</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/sprout-cloudflared.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sprout-cloudflared.err.log</string>
</dict>
</plist>
