# infra/launchd/com.sprout.activities.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly activity generation. Sunday 6:00am, before the digest (6:30) so the
  digest can include this week's activities.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sprout.activities</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/pnpm</string>
    <string>tsx</string>
    <string>scripts/weekly-activities.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/waylantrieu/repos/sprout</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>DATABASE_URL</key>
    <string>/Users/waylantrieu/repos/sprout/data/sprout.db</string>
    <key>OLLAMA_URL</key>
    <string>http://localhost:11434</string>
    <key>OLLAMA_MODEL</key>
    <string>qwen2.5:14b</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/sprout-activities.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sprout-activities.err.log</string>
</dict>
</plist>
